# systemd/network-check-web.service
[Unit]
Description=Network Check web dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__WORKDIR__
ExecStart=__NODE_BIN__ server/index.js
Restart=on-failure
RestartSec=5
Environment=PORT=3000
Environment=HOST=0.0.0.0

[Install]
WantedBy=default.target
