# systemd/network-check-collect.service
[Unit]
Description=Collect one network speedtest sample
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__WORKDIR__
ExecStart=__NODE_BIN__ server/collect.js
Nice=10
